# summer2014/allskycam/allskycam.crontab
MAILTO=""

*/5 * * * * /data/allskycam/allskycam_fiveminutes.sh
